# root/etc/s6-overlay/s6-rc.d/init-workdir/up
#!/command/execlineb -P
with-contenv
define main_dir "/epgd"
define conf_dir "${main_dir}/config"
foreground {
 if -n { eltest -d "${conf_dir}" }
  s6-mkdir -p -m 0755 "${conf_dir}"
 }
foreground { s6-echo "Prepare workdir" }
cp -RPpn /defaults/config/. "${conf_dir}"
