# root/etc/services.d/epghttpd/run
#!/usr/bin/with-contenv bash

if [ "$START_EPGHTTPD" = "yes" ]; then
	echo "Starting epghttpd..."
	export LANG="$LANG"
	exec \
		s6-setuidgid abc /usr/local/bin/epghttpd -n -t -c /epgd/config
fi
