# root/etc/services.d/epghttpd/run
#!/usr/bin/with-contenv bash

if [ "$START_EPGHTTPD" = "yes" ]; then
	s6-echo "Starting epghttpd..."
	exec \
		s6-setuidgid abc /usr/local/bin/epghttpd -n -c /epgd/config
fi
