# root/etc/s6-overlay/s6-rc.d/restart-epgd/up
#!/command/execlineb -P
with-contenv
backtick -D "false" restart { printcontenv RESTART_EPGD }
backtick -D "3:30" time { printcontenv RESTART_EPGD_AT }
multisubstitute
{
  importas -usC restart restart
  importas -usC time time
}
if { test $restart != false }
 if -n { test -z $time }
  foreground { s6-echo "Planing epgd restart..." }
  pipeline { s6-echo "s6-svc -r /run/service/epgd" }
  /usr/bin/at $time
