# root/etc/s6-overlay/s6-rc.d/epghttpd/run
#!/command/execlineb -P
with-contenv
importas -D "true" start START_EPGHTTPD
ifelse { s6-test $start = false }
 {
  foreground { s6-echo "epghttpd is disabled by environment variable!" }
  foreground { s6-rmrf /etc/s6-overlay/s6-rc.d/user/contents.d/epghttpd }
  s6-svc -d -O /run/s6-rc/servicedirs/epghttpd
 }
 s6-setuidgid epgd
 fdmove -c 2 1
 /usr/local/bin/epghttpd -n -c /epgd/config
