# root/etc/services.d/epgd/run
#!/usr/bin/with-contenv bash

s6-echo "Starting epgd..."
exec \
	s6-setuidgid abc /usr/local/bin/epgd -n -p /usr/local/lib/epgd/plugins -c /epgd/config
