# root/etc/services.d/epgd/run
#!/usr/bin/with-contenv bash

echo "Starting epgd..."
export LANG="$LANG"
exec \
	s6-setuidgid abc /usr/local/bin/epgd -n -t -p /usr/local/lib/epgd/plugins -c /epgd/config
