# root/etc/s6-overlay/s6-rc.d/epgd/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
s6-setuidgid epgd
/usr/local/bin/epgd -n -p /usr/local/lib/epgd/plugins -c /epgd/config
